# Daily_crawler/daily_crawler.cron
0 23 * * * /home/johnso/PycharmProjects/News_recommendation/News_scrapy_redis/Daily_crawler/start_crawl.sh
